# wifi-button/gestures.service
[Unit]
Description=Pushbutton+LED gestures daemon
After=wifi-init.service
Requires=wifi-init.service

[Service]
Type=exec
WorkingDirectory=/opt/sensorgnome/wifi-button/scripts
ExecStartPre=/usr/bin/bash gpio-init.sh ../misc/gestures.txt
ExecStart=/usr/bin/node gestures.js ../misc/gestures.txt
SyslogIdentifier=sg-gestures
RestartSec=60s
Restart=on-failure

[Install]
WantedBy=multi-user.target
