# cellular/check-modem.timer
[Unit]
Description=Check modem timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=multi-user.target
